# scripts/files.req.list
/etc/X11/profile.d
/etc/X11/wmsession.d
/etc/X11/xinit.d
/etc/alternatives/packages.d
/etc/chroot.d
/etc/cron.d
/etc/cron.daily
/etc/cron.hourly
/etc/cron.monthly
/etc/cron.weekly
/etc/emacs/site-start.d
/etc/logrotate.d
/etc/pam.d
/etc/ppp/ip-down.d
/etc/ppp/ip-up.d
/etc/rpm/macros.d
/etc/security/console.apps
/etc/sudo.d
/etc/syslog.d
/etc/tex-fonts.d
/etc/udev/rules.d
/etc/xinetd.d
